# debian/vyatta-cfg-system.postrm
#!/bin/bash

set -e

if [ "$1" = "purge" ]; then
	systemctl disable vyatta-config-reboot-params.service
	systemctl disable vyos-intfwatchd.service
	systemctl disable ec2-fetch-ssh-public-key.service
fi

exit 0
